# bin/systemd/meta3/gunicorn-live.socket
[Unit]
Description=gunicorn live socket

[Socket]
ListenStream=/run/gunicorn-live.sock

[Install]
WantedBy=sockets.target
