# bin/systemd/gunicorn.service
# confirm the working directory

[Unit]
Description=gunicorn daemon
Requires=gunicorn.socket
After=network.target

[Service]
User=live
Group=www-data
WorkingDirectory=/home/live/app
ExecStart=/home/live/.venvs/edc/bin/gunicorn \
    --access-logfile - \
    --workers 3 \
    --timeout=300 \
    --bind unix:/run/gunicorn.sock \
    meta_edc.wsgi:application

[Install]
WantedBy=multi-user.target